# timeular/.config/systemd/user/timeular.service
[Unit]
Description=Timeular service
After=bluetooth.target

[Service]
ExecStart=/usr/bin/timeular

[Install]
WantedBy=multi-user.target
